# environment/ansible/labs/D04/ping.service
[Unit] 
Description="Write ping to log"

[Service]
ExecStart=/usr/local/bin/ping.sh
