# environment/ansible/labs/D04/ping.timer
[Unit]
Description=Ping every minute

[Timer]
OnCalendar=*:0/1
Unit=ping.service

[Install]
WantedBy=timers.target
