# environment/ansible/labs/D03/passwd.path
[Unit]
Description="Monitor /etc/passwd for changes"

[Path]
PathModified=/etc/passwd
Unit=passwd-log.service

[Install]
WantedBy=multi-user.target
