# environment/ansible/labs/D03/passwd-log.service
[Unit] 
Description="Write alarm to passwd log"

[Service]
ExecStart=/usr/local/bin/passwd-alert.sh
