# environment/ansible/labs/D03/ramdisk.mount
[Unit]
Description=Mount smol ramdisk

[Mount]
What=none
Where=/ramdisk
Type=tmpfs
Options=size=200M

[Install]
WantedBy=multi-user.target
